# runit/.config/runit/sv/xautolock/run
#!/bin/sh
xautolock -time 15 -locker slock -detectsleep -killtime 10 -killer zzz -resetsaver
