# runit/.config/runit/sv/emacs/run
#!/bin/sh
cd $HOME
exec emacs --fg-daemon
