# runit/.config/runit/sv/mpd/run
#!/bin/sh

[ -r ./conf ] && . ./conf

exec mpd --no-daemon --stderr ${OPTS:-} 2>&1
